# app/installer/config/supervisor/newapp_staging.conf
[program:newapp-staging]
command=/usr/bin/node /home/ubuntu/newapp-staging/src/api/index.js
process_name=%(program_name)s
numprocs=1
directory=/home/ubuntu/newapp-staging
umask=022
priority=999
autostart=false
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=false
environment=
  NODE_ENV=staging,
  NODE_PORT=8080
stdout_logfile=/var/log/newapp_staging.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/newapp_staging.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
serverurl=AUTO
